# custom-script
#! /bin/sh
# chkconfig: 345 99 10
case "$1" in
  start)
    # Executes our script
    sudo sh ./start.sh
    ;;
  *)
    ;;
esac
exit 0
